# scripts/com.milkroad.onchain-index-docs-serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.milkroad.onchain-index-docs-serve</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/max/projects/onchain-index/.venv/bin/python</string>
        <string>-m</string>
        <string>http.server</string>
        <string>8012</string>
        <string>--bind</string>
        <string>0.0.0.0</string>
        <string>--directory</string>
        <string>/Users/max/projects/onchain-index/docs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/max/projects/onchain-index</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/max/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/max/projects/onchain-index/.cache/launchd-docs-serve.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/max/projects/onchain-index/.cache/launchd-docs-serve.log</string>
</dict>
</plist>
